# systemd/freeipa-backup.timer
[Unit]
Description=FreeIPA Backup Timer
Documentation=file:///usr/local/share/freeipa-backup-automation/README.md
Requires=freeipa-backup.service

[Timer]
# Run daily at 2:00 AM
OnCalendar=*-*-* 02:00:00

# Run 15 minutes after boot if the system was off during the scheduled time
OnBootSec=15min

# Allow up to 30 minutes of randomized delay to avoid thundering herd
RandomizedDelaySec=30min

# If the backup takes longer than expected, don't start another one
Persistent=true

# Wake up the system if it's sleeping (requires appropriate power management setup)
WakeSystem=true

[Install]
WantedBy=timers.target
